# system/etc/init/koule.conf
#!upstart
#
description "koule for upstart"

respawn

start on tty-device-added 
stop on runlevel [016]

exec /home/odroid/git/balon/scripts/m_koule.py

